# socketTest.socket
[Socket]
ListenStream=/run/socketTest.sk

[Install]
WantedBy=sockets.target